# roles/goss/templates/prometheus_goss.cron.j2
*/1 * * * * root {{ goss_install_dir }}/prometheus_goss.sh
